# ansible/templates/com.miniai.gateway.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
{# Rendered by ansible/site.yml — do not edit the installed copy by hand.
   Note: NO configuration here. The gateway pulls everything from SSM
   Parameter Store at startup (server/config.py, terraform/ssm.tf). #}
<plist version="1.0">
<dict>
  <key>Label</key><string>com.miniai.gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{ python_bin }}</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>server.gateway:app</string>
    <string>--host</string>
    <!-- 0.0.0.0 so the bastion can reach us over the WireGuard interface
         (10.8.0.2). The mlx model server stays on 127.0.0.1 — only the
         gateway is tunnel-facing. Keep the macOS firewall on. -->
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8000</string>
  </array>
  <key>WorkingDirectory</key><string>{{ miniai_dir }}</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- profile NAME is not a secret; the credential lives in ~/.aws.
         This lets server/config.py reach SSM when running under launchd. -->
    <key>AWS_PROFILE</key><string>miniai</string>
    <key>AWS_REGION</key><string>us-east-1</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>{{ miniai_dir }}/logs/gateway.log</string>
  <key>StandardErrorPath</key><string>{{ miniai_dir }}/logs/gateway.log</string>
</dict>
</plist>
